# docker/collectd-supervisor.conf
[program:collectd]
directory=/usr/lib/
command=/etc/init.d/collectd restart
autostart=true
autorestart=true
redirect_stderr=true
